# deploy/expo-go/luche-go-expo.service
[Unit]
Description=Luche Expo Go Metro server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi-rus/Downloads/feral-remote/luche-go/current
Environment=CI=1
Environment=EXPO_NO_TELEMETRY=1
Environment=EXPO_PACKAGER_PROXY_URL=https://go.luche.ai
Environment=NODE_OPTIONS=--max-old-space-size=3072
ExecStart=/usr/bin/node node_modules/expo/bin/cli start --go --lan --port 8091 --no-dev --minify --max-workers 2 --clear
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=default.target
